# deployments/envoy-agent.service
[Unit]
Description=Envoy VM Agent
After=network-online.target
Wants=network-online.target

[Service]
Environment=ENVOY_MANIFEST=/etc/envoy/agent.manifest.yaml
Environment=ENVOY_NATS_URL=nats://platform.example.com:4222
ExecStart=/usr/local/bin/envoy-agent
Restart=always
RestartSec=5
User=envoy

[Install]
WantedBy=multi-user.target
